# package/pwst/radioconf/files/etc/hotplug.d/ieee80211/99-morse-duty-cycle
#!/bin/sh

# Extract the phy name
PHY="${DEVPATH##*/ieee80211/}"

# Only process phy0 (the Morse Micro device on spi0.0 = radio0)
[ "$PHY" = "phy0" ] || exit 0

# Log for debugging
logger -t morse-duty-cycle "Morse radio0 (phy0) event detected: $ACTION"

# Wait a moment for the radio to be fully initialized
sleep 3

# Execute your duty cycle configuration script
/usr/bin/healthchek-n-conf.sh

logger -t morse-duty-cycle "Duty cycle settings applied to radio0"
