# target/linux/ramips/mt7621/base-files/etc/hotplug.d/leds/10-pwm-multicolor
#!/bin/sh

# OpenWrt has no mechanism yet to set a multicolor LEDs RGB colour value. For now, set the colour value to
# a purple on insertion of the leds-pwm-multicolor driver.
# For intial demonstration of the PWM RGB LEDs. To be improved later.

if [ "$ACTION" = add ]; then
	basename "$(readlink -f "/sys$DEVPATH/device/driver/")" | grep '^leds_pwm_multicolor' || return 0;
	COLOR="255 0 100"
	expr $DEVPATH : ".*led0.*" && COLOR="0 255 0"
	echo $COLOR > "/sys$DEVPATH/multi_intensity"
fi
